# deploy/filesearch.service
[Unit]
Description=File Search RAG Backend
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/filesearch
Environment="PATH=/var/www/filesearch/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/filesearch/venv/bin/gunicorn --worker-class eventlet -w 1 --bind 127.0.0.1:5555 "backend.app:create_app()"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
